# crontab
* * * * * bash ~/scripts/daemon.sh &> /dev/null
